# ansible/roles/openbao/templates/openbao-snapshot.timer.j2
{# Managed by the openbao role. Persistent so a missed run (host down at 03:00)
   fires on next boot. #}
[Unit]
Description=Daily OpenBao raft snapshot to S3

[Timer]
OnCalendar={{ openbao_snapshot_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
